# recipes-core/systemd/systemd-genivi-targets/focussed.target
[Unit]
Description=Focussed Target
After=node-state-controller.service
Requires=basic.target
AllowIsolate=yes
