# recipes-core/systemd/systemd-genivi-targets/lazy.target
[Unit]
Description=Lazy Target
Requires=unfocussed.target multi-user.target
After=unfocussed.target
AllowIsolate=yes
