# recipes-core/systemd/systemd-genivi-targets/unfocussed.target
[Unit]
Description=Unfocussed Target
Requires=focussed.target
After=focussed.target
AllowIsolate=yes
