# recipes-core/systemd/systemd-additional-units/light.target
[Unit]
Description=Light Target
Requires=dbus.socket
